# crontab
*/1 * * * * echo $(data)>>/var/log/cron.log
